# pi-backend/router-api.service
[Unit]
Description=Travel Router API
After=network.target
Wants=network.target

[Service]
Type=simple
User=router-api
Group=router-api
WorkingDirectory=/opt/travel-router-admin/pi-backend
Environment=PATH=/opt/travel-router-admin/pi-backend/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/opt/travel-router-admin/pi-backend/venv/bin/python -m uvicorn src.main:app --host 0.0.0.0 --port 8080
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
